# zen-claw.service
[Unit]
Description=Zen Claw Gateway
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
Restart=always
RestartSec=1
User=neves
ExecStart=/home/neves/zen/zen-claw/zen-claw gateway start
ExecStop=/home/neves/zen/zen-claw/zen-claw gateway stop
WorkingDirectory=/home/neves/.zen/zen-claw
Environment="PATH=/home/neves/go/bin:/usr/local/go/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=zen-claw

# Security (relaxed to allow home directory access for config)
NoNewPrivileges=true
ProtectSystem=strict
PrivateTmp=true
PrivateDevices=true
ProtectHome=read-only
ReadWritePaths=/tmp /home/neves/.zen
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target